# systemd/zfs-scrub@.service
[Unit]
Description=Scrub %I

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/sbin/zpool scrub %I
