# systemd/zfs-snapshot-monthly@.timer
[Unit]
Description=Monthly %I snapshots

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
