# systemd/zfs-snapshot-weekly@.timer
[Unit]
Description=Weekly %I snapshots

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
